# packages/borgee-helper/install/pkg/scripts/postinstall
#!/bin/sh
# macOS postinstall — runs after payload is in place.
#
# Bootstraps the LaunchDaemon so it starts on next boot but does NOT start
# it now: the heartbeat producer needs claim files first (same rationale as
# the Linux postinstall.sh).

set -eu

PLIST="/Library/LaunchDaemons/cloud.borgee.host-bridge.plist"

if [ -f "$PLIST" ]; then
    chown root:wheel "$PLIST"
    chmod 0644 "$PLIST"
    if command -v launchctl >/dev/null 2>&1; then
        launchctl bootstrap system "$PLIST" 2>/dev/null || true
        # Intentionally not `launchctl kickstart` — wait for operator claim
        # before the daemon should actually start producing heartbeats.
    fi
fi

cat <<'EOF'

Borgee Helper installed.

Next steps:
  1. Generate an enrollment in the Borgee web UI (Hosts -> Add).
  2. On this Mac, run:

       sudo borgee-helper-claim \
           --enrollment-id   <ID from web UI> \
           --enrollment-secret <one-time secret>

  3. Start the daemon:

       sudo launchctl kickstart -k system/cloud.borgee.host-bridge

  4. Tail logs:

       tail -f /Library/Logs/Borgee/stderr.log
EOF

exit 0
